# mtk-openwrt-4.0.1.0/files/etc/init.d/netmgrd
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/sbin/netmgrd &
	echo "run netmgrd.."
}

stop() {
	echo "stop gunicorn..."
	killall netmgrd 
}
