# mtk-openwrt-4.0.1.0/package/mtk/misc/autotest/files/autotest.init
#!/bin/sh /etc/rc.common

START=99

USE_PROCD=1
PROG=/usr/sbin/autotest

start_service() {

	. /usr/share/autotest/autotest.sh

}


